# launchd/memory-distill.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Template for the per-agent session-distillation launchd job.
  The install scripts (scripts/install_codex_watcher.sh,
  scripts/install_cursor_watcher.sh) substitute the __PLACEHOLDER__ tokens
  with machine-specific paths at install time and write the result to
  ~/Library/LaunchAgents/. This template ships with placeholders only, so no
  absolute paths are baked into the repo.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>__LABEL__</string>

  <key>ProgramArguments</key>
  <array>
    <string>__PYTHON__</string>
    <string>__SCANNER__</string>
    <string>--quiet-minutes</string>
    <string>30</string>
    <string>--lookback-days</string>
    <string>3650</string>
    <string>--limit</string>
    <string>0</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>__PATH__</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>4</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>__LOG_OUT__</string>

  <key>StandardErrorPath</key>
  <string>__LOG_ERR__</string>
</dict>
</plist>
